# root/etc/s6-overlay/s6-rc.d/svc-ledfx/run
#!/command/with-contenv bash
# Start LedFx using s6
exec /usr/bin/ledfx
